# every-15m.cron
*/15 8-21 * * * /home/supremine/supremine
